# rpc-server.service
[Unit]
Description=ZeroRpc Server

[Service]
Type=simple
ExecStart=/bin/bash /home/root/usrp_uhd_api/start_usrp_server.bash
ProtectHome=off

[Install]
WantedBy=multi-user.target
